# rpi-logger.service
[Unit]
Description=RPi Logger - Multi-module data logging system
Documentation=file:///home/rs-pi-2/Development/RPi_Logger/CLI_MODE.md
After=graphical.target network.target
Wants=graphical.target

[Service]
Type=simple
User=rs-pi-2
Group=rs-pi-2
WorkingDirectory=/home/rs-pi-2/Development/RPi_Logger

# Environment variables
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/rs-pi-2/.Xauthority"
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Startup script (reads config.txt for mode)
ExecStart=/home/rs-pi-2/Development/RPi_Logger/start_logger.sh

# Restart on failure (optional - comment out if you don't want auto-restart)
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rpi-logger

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
